# debbuilder/nubomanagement-serv-common/DEBIAN/prerm
#!/bin/sh

if [ "$1" = remove ]; then
  systemctl disable nubomanagement
fi
systemctl stop nubomanagement
