# docker_build/etc/supervisor/conf.d/nubomanagement.conf
[program:nubomanagement]
directory=/opt/nubomanagement
command=/usr/bin/node dist/multithreadserver.js
user=root
autorestart=true
priority=400
